# patterns/03_immune-loop/immune-detector.timer
[Unit]
Description=Immune Loop — Daily Violation Scan
Requires=immune-detector.service

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
